# config/screenlayout/screenlayout.rules
SUBSYSTEM=="drm", ACTION=="change", RUN+="/home/$USER/.screenlayout.sh"
